# examples/etc/rc.d/xbsapi
#!/bin/ksh

daemon="/usr/local/bin/xbsapi"

. /etc/rc.d/rc.subr

rc_start() {
  ${rcexec} "${daemon} ${daemon_flags} 2>&1 | logger -t xbsapi &"
}

rc_cmd $1
